# templates/luigi.conf
description "Luigi"
author  "Space/Time Directory <spacetime@nypl.org>"

start on filesystem and started networking

exec luigid
respawn
